# docker/root-fs/etc/sv/nginx/run
#!/bin/sh

##
# Runit run script for nginx
#

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
